# systemd/runner-workspace-cleanup.timer
[Unit]
Description=Nightly GitHub Actions workspace cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
